# mirror/provisioning/bluetooth.service.d/10-smartmirror.conf
# Installed to /etc/systemd/system/bluetooth.service.d/ by provisioning/install.sh.
#
# Beats the Pi 5 boot-race: bluetooth.service has
# ConditionPathIsDirectory=/sys/class/bluetooth, which condition-skips bluetoothd
# if it's evaluated before the kernel has attached the BCM4345C0 controller (hci0).
# Ordering after the hci0 device makes the condition hold, and Restart recovers
# bluetoothd if it ever exits while still needed.
[Unit]
After=sys-subsystem-bluetooth-devices-hci0.device
Wants=sys-subsystem-bluetooth-devices-hci0.device

[Service]
Restart=on-failure
RestartSec=2
