# mirror/provisioning/smartmirror-sync.service
[Unit]
Description=Smart Mirror Sync Bridge (:4002)
Documentation=https://github.com/YassoBases/Smart-Mirror-Combined
After=network.target smartmirror-backend.service
Wants=smartmirror-backend.service

[Service]
Type=simple
# __MIRROR_DIR__ and __NODE_BIN__ are rewritten to the real deploy location by install-app.sh.
ExecStart=__NODE_BIN__ sync/dist/server.js
WorkingDirectory=__MIRROR_DIR__
User=smartmirror
Environment=MIRROR_BACKEND_URL=ws://localhost:4000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
